# deploy/systemd/newdomofon-video-master-disk-guard.service
[Unit]
Description=NewDomofon Video Master Disk Guard
Documentation=file:/opt/newdomofon-video-master/docs/DISK_PROTECTION.md
After=local-fs.target postgresql.service

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=-/etc/newdomofon-video/app.env
ExecStart=/usr/local/sbin/newdomofon-master-disk-guard
Nice=10
IOSchedulingClass=idle
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=false
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/run/newdomofon-video /run/lock /var/log/newdomofon-video /tmp /var/tmp /var/cache/apt
CapabilityBoundingSet=
AmbientCapabilities=
